# scripts/prompt-manager.service
[Unit]
Description=Prompt Manager Application
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=notify
User=www-data
Group=www-data
WorkingDirectory=/var/www/prompt-manager
Environment="PATH=/var/www/prompt-manager/venv/bin"
Environment="FLASK_ENV=production"
ExecStart=/var/www/prompt-manager/venv/bin/gunicorn \
    --config /var/www/prompt-manager/gunicorn.conf.py \
    --access-logfile /var/log/prompt-manager/access.log \
    --error-logfile /var/log/prompt-manager/error.log \
    wsgi:app

# Restart policy
Restart=always
RestartSec=5

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/www/prompt-manager /var/log/prompt-manager

# Resource limits
LimitNOFILE=65535
LimitNPROC=4096

# Logging
StandardOutput=append:/var/log/prompt-manager/service.log
StandardError=append:/var/log/prompt-manager/service-error.log

[Install]
WantedBy=multi-user.target